# network-manager-customsettings-cups/etc/systemd/system/update.target
[Unit]
Description=Daily update Target
#StopWhenUnneeded=yes

