# network-manager-customsettings/usr/lib/systemd/system/update.target
[Unit]
Description=Target for running several updates
StopWhenUnneeded=no
